# upstart.conf.template
description "uWSGI for {UWSGI_FOR}"
start on runlevel [2345]
stop on runlevel [06]

exec uwsgi --master --die-on-term --ini {UWSGI_INI_LOCATION}
